# version1.1/etc/cron.hourly/logrotate
#!/bin/sh

sudo /usr/sbin/logrotate /etc/logrotate.conf
